# dot_config/shell/myscript.service
[Unit]
Description=My script
[Service]
ExecStart=/usr/bin/my-script
Type=idle
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
